# systemd/user/clean.timer
[Unit]
Description=Pulizia

[Timer]
OnCalendar=6:00
Persistent=true     

[Install]
WantedBy=timers.target
