# systemd/user/clean.service
[Unit]
Description=Pulizia

[Service]
Type=oneshot
ExecStart=%h/.homeassistant/clean.sh
